# testconfigs/shared/overlay-livecd/etc/local.d/rc.local.start
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#
/usr/bin/rc-service checknet start  # Check if networkmanager started correctly
/usr/bin/pacman-boot &  # Setup pacman databases
# Source /etc/conf.d/xdm if it exists
[ -e /etc/conf.d/xdm ] && . /etc/conf.d/xdm
if [ "$DISPLAYMANAGER" == "lxdm" ]; then
	rc-service xdm restart  # Restart xdm as lxdm does not get started correctly the first time, does not happen with lightdm
fi
