# testconfigs/lxqt-openrc/lxqt-overlay/etc/skel/.config/autostart/Connman.desktop
[Desktop Entry]
Exec=cmst -m -w 3
Name=Connman
Icon=cmst
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
